# WorkData/tools/tagging/data_consistency/daily_checks.sh.crontab
MAILTO=straka@ufal.mff.cuni.cz

# m h  dom mon dow   command
0   4  *   *   *     . ~/.profile && cd /home/straka/repos/pdt-c/WorkData/tools/tagging/data_consistency/ && sh daily_checks.sh
